# .config/systemd/user/telegram.service
[Unit]
Description=Telegram

[Service]
ExecStart=%h/bin/Telegram/Telegram
Restart=on-abort

[Install]
WantedBy=graphical.target
