# .config/systemd/user/udiskie.service
[Unit]
Description=Udiskie

[Service]
ExecStart=/usr/bin/udiskie -s
Restart=on-abort

[Install]
WantedBy=graphical.target
