# .config/systemd/user/zeal.service
[Unit]
Description=Zeal

[Service]
ExecStart=/usr/bin/zeal
Restart=on-abort

[Install]
WantedBy=graphical.target
